# debian/addrnodeimport.cron.d


10 2 * * 7 addrnodeimport /usr/lib/addrnodeimport/sh/trigger_update.sh > /dev/null 2>&1

